# supervisord.conf
[supervisord]
nodaemon=true

[program:flask]
command=gunicorn app:app --bind 0.0.0.0:5000

[program:streamlit]
command=streamlit run dashb.py --server.port 8501 --server.address 0.0.0.0
